# systemd/patient-data-collection.service
[Unit]
Description=Patient Data Collection API
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/patient-data-collection
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=append:/var/log/patient-data-collection/access.log
StandardError=append:/var/log/patient-data-collection/error.log

[Install]
WantedBy=multi-user.target